# taginfo.ini
[program:taginfo]
command=/usr/local/bin/ruby taginfo.rb 2002
autostart=true
autorestart=true
stopsignal=INT
stopwaitsecs=30
user=nobody
directory=/taginfo-run/repo/web
